# docker/crontab
* * * * * preapp_dev /usr/local/bin/php /var/www/test/artisan schedule:run >> /dev/null 2>&1